# systemd/user/rustserver-lgsm-update.timer
[Unit]
Description=rustserver lgsm-update timer

[Timer]
OnUnitActiveSec=5m
OnBootSec=10s

[Install]
WantedBy=timers.target
